# docker/supervisor/conf.d/stream.conf
[supervisord]
nodaemon=true

[program:stream]

command = gunicorn "app:create_app()" --bind=0.0.0.0:8021 --timeout=180 
directory = /workspace
autostart=true
autorestart=true
